# scripts/calendars.crontab
# Minute   Hour   Day of Month       Month          Day of Week        Command    
# (0-59)  (0-23)     (1-31)    (1-12 or Jan-Dec)  (0-6 or Sun-Sat)                
#
#
#
# Description:  beavo crontab update dcppc calendar script
# Frequency:    hourly
#
#   MIN      HR         DOM            MONTH            DOW          CMD
    0        *          *              *                *            /home/ubuntu/calendars/scripts/update_calendars.sh
